# deploy/systemd/headless-scraping-subsystem.service
[Unit]
Description=Headless Scraping Subsystem (Imperva challenge solver)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=headless
Group=headless
WorkingDirectory=/opt/headless-scraping-subsystem
EnvironmentFile=/etc/headless-scraping-subsystem.env
ExecStart=/opt/headless-scraping-subsystem/.venv/bin/uvicorn app.main:app \
          --host ${BIND_HOST} --port ${BIND_PORT} \
          --workers 1 --log-level ${LOG_LEVEL}
Restart=on-failure
RestartSec=5

# Hardening
# NB: Chromium uses AF_NETLINK for network change detection and unshare()
# for its zygote sandbox, so we cannot enable LockPersonality and we must
# allow AF_NETLINK in RestrictAddressFamilies. The remaining hardening is
# still meaningful: no new privileges, read-only system, private /tmp.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictRealtime=true
MemoryHigh=2G
MemoryMax=3G
TasksMax=512

# Playwright places its browser cache here at install time.
ReadWritePaths=/opt/headless-scraping-subsystem

[Install]
WantedBy=multi-user.target
